# triode.service.in
[Unit]
Description=Triode
After=network.target

[Service]
WorkingDirectory=@SYSCONFDIR@/@ASSEMBLY_NAME@
ExecStart=@SBINDIR@/@ASSEMBLY_NAME@
Restart=always
RestartSec=10
SyslogIdentifier=@ASSEMBLY_NAME@
User=@SERVICE_USER@
AmbientCapabilities=CAP_NET_BIND_SERVICE
Environment=ASPNETCORE_ENVIRONMENT=Production
Environment=DOTNET_GENERATE_ASPNET_CERTIFICATE=false
Environment=HOME=@LOCALSTATEDIR@/lib/@ASSEMBLY_NAME@

[Install]
WantedBy=multi-user.target
